# crontab
0 0 1 * * node /home/runcloud/webapps/Automation/app.js